# packages/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:backend]
command=node /app/backend/index.js
directory=/app/backend
environment=NODE_ENV=%(ENV_NODE_ENV)s,PORT=%(ENV_PORT)s,ORACLE_SERVER=%(ENV_ORACLE_SERVER)s,ORACLE_PORT=%(ENV_ORACLE_PORT)s,ORACLE_USERNAME=%(ENV_ORACLE_USERNAME)s,ORACLE_PASSWORD=%(ENV_ORACLE_PASSWORD)s,ORACLE_DATABASE=%(ENV_ORACLE_DATABASE)s,JWT_SECRET=%(ENV_JWT_SECRET)s
autostart=true
autorestart=true
stderr_logfile=/var/log/app/backend.err.log
stdout_logfile=/var/log/app/backend.out.log
user=node
startsecs=10
startretries=3

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
user=root
depends_on=backend